# package/mypackages/luci-app-chinaroutes/files/root/etc/uci-defaults/luci-chnroute
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@chnroute[-1]
	add ucitrack chnroute
	set ucitrack.@chnroute[-1].init=chnroute
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
