# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

ipset create ss-cmcc hash:ip timeout 3600
ipset create ss-telecom hash:ip timeout 3600

exit 0
